# rose-receptionist.service
[Unit]
Description=Rose Receptionist API Service
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/opt/rose-receptionist
ExecStart=/usr/bin/dotnet /opt/rose-receptionist/RoseReceptionist.API.dll
Restart=always
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=rose-receptionist
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://localhost:5000

[Install]
WantedBy=multi-user.target
